# config/com.rapid-mlx.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rapid-mlx.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mkoziol/.local/bin/rapid-mlx</string>
        <string>serve</string>
        <string>qwen3.6-35b</string>
        <string>--port</string>
        <string>8000</string>
        <string>--no-mllm</string>
        <string>--max-num-seqs</string>
        <string>1</string>
        <string>--max-concurrent-requests</string>
        <string>4</string>
        <string>--reasoning-parser</string>
        <string>qwen3</string>
        <string>--default-temperature</string>
        <string>0.6</string>
        <string>--default-top-p</string>
        <string>0.95</string>
        <string>--default-top-k</string>
        <string>20</string>
        <string>--max-tokens</string>
        <string>32768</string>
        <string>--stream-interval</string>
        <string>1</string>
        <string>--pin-system-prompt</string>
        <string>--enable-prefix-cache</string>
        <string>--kv-cache-turboquant</string>
        <string>--kv-cache-turboquant-bits</string>
        <string>3</string>
        <string>--cache-memory-mb</string>
        <string>2000</string>
        <string>--gpu-memory-utilization</string>
        <string>0.9</string>
        <string>--prefill-step-size</string>
        <string>4096</string>
        <string>--chunked-prefill-tokens</string>
        <string>4096</string>
        <string>--gc-control</string>
        <string>--enable-auto-tool-choice</string>
        <string>--tool-call-parser</string>
        <string>qwen3_coder_xml</string>
        <string>--timeout</string>
        <string>1800</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RAPID_MLX_TELEMETRY</key>
        <string>0</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/rapid-mlx.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rapid-mlx.err</string>
</dict>
</plist>
